# rootfs/etc/systemd/system/etc/rc.local
[ -d /etc/boot.d ] && run-parts /etc/boot.d

pigpiod

exit 0
